# infra/systemd/wb-helper-compose.service
[Unit]
Description=WB Helper Docker Compose Stack
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/wb_helper
Environment=COMPOSE_FILE=infra/docker-compose.yml
EnvironmentFile=-/etc/wb-helper/wb-helper.env
ExecStart=/usr/bin/docker compose up -d --remove-orphans
ExecStop=/usr/bin/docker compose down
ExecReload=/usr/bin/docker compose up -d --remove-orphans
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
